# templates/default/awslogs.service.erb
Description=The CloudWatch Logs Agent
After=rc-local.service

[Service]
Type=simple
Restart=always
KillMode=process
TimeoutSec=infinity
Nice=4
PIDFile=<%= node['cwlogs']['base_dir'] %>/state/awslogs.pid
EnvironmentFile=-<%= node['cwlogs']['base_dir'] %>/etc/proxy.conf
EnvironmentFile=-/etc/default/awslogs
Environment="HOME=<%= node['cwlogs']['base_dir'] %>"
Environment="AWS_CONFIG_FILE=<%= node['cwlogs']['base_dir'] %>/etc/aws.conf"
ExecStart=<%= node['cwlogs']['base_dir'] %>/bin/aws logs push --config-file <%= node['cwlogs']['base_dir'] %>/etc/awslogs.conf --additional-configs-dir <%= node['cwlogs']['base_dir'] %>/etc/config

[Install]
WantedBy=multi-user.target
